# vagrant/cookbooks/graphite/templates/default/supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=<%= node[:graphite][:venv] %>/run/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=<%= node[:graphite][:venv] %>/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:carbon]
; nondaemonized startup method required
; http://supervisord.org/subprocess.html
command=/opt/graphite/bin/carbon-cache.py --debug start

[program:graphite-devel-server]
; nondaemonized startup method required
command=/opt/graphite/bin/run-graphite-devel-server.py /opt/graphite/webapp/graphite

[program:gunicorn]
; nondaemonized startup method required
; http://supervisord.org/subprocess.html
command=<%= node[:graphite][:venv] %>/bin/gunicorn_django --preload -w 2 -b 127.0.0.1 --log-level debug --log-file <%= node[:graphite][:venv] %>/run/gunicorn.log -p <%= node[:graphite][:venv] %>/run/gunicorn.pid /opt/graphite/webapp/graphite/settings.py
redirect_stderr=true

[program:node]
; nondaemonized startup method required
; http://supervisord.org/subprocess.html
command=/usr/local/bin/node /opt/statsd/stats.js /etc/statsd.js

[group:graphite]
; common parameters
programs=node,carbon,gunicorn,graphite-devel-server
numprocs=1
startsecs=3
autostart=true
autorestart=true
startretries=3
exitcodes=0,2
